# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:system_info]
command=/usr/share/nginx/html/system_info.sh
